# SourceCode/Script.cron
# Fetch Stock Price
# Run at 8:59AM and 12:59PM from Mon to Fri
59 1,5 * * 1-5 /usr/bin/python3 ~/StockFetchPrice/fetch_stock_price.py
# Run from 9:00 to 10:59 and from 13:00 to 14:59 from Mon to Fri
0-59 2-3,6-7 * * 1-5 /usr/bin/python3 ~/StockFetchPrice/fetch_stock_price.py
# Run from 11:00 to 11:34 from Mon to Fri
0-34 4 * * 1-5 /usr/bin/python3 ~/StockFetchPrice/fetch_stock_price.py
# Run from 15:00 to 15:05 from Mon to Fri
0-5 8 * * 1-5 /usr/bin/python3 ~/StockFetchPrice/fetch_stock_price.py
# --------------------------------------------------------------------
# Query for prediction
# Run at 3:30PM from Mon to Fri
30 8 * * 1-5 /usr/bin/python3 ~/QueryForPrediction/query_for_prediction.py
#---------------------------------------------------------------------
# Company Info 
# Run at 4:00 PM from Mon to Fri
0 9 * * 1-5 /usr/bin/python3 ~/CompanyInfo/fetch_company_info.py
